# templates/supervisord.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn welgo:app -c {{ code_dir }}/gunicorn.conf.py
directory={{ code_dir }}
user=root
stdout_logfile=/var/log/gunicorn-access.log
stderr_logfile=/var/log/gunicorn-error.log
autostart=true
autorestart=true
environment=PATH="/usr/local/bin:/usr/bin:/bin", PYTHONUNBUFFERED=true


[program:celery]
command=/usr/local/bin/celery worker -A welgo:celery --concurrency=3 --loglevel=INFO

directory={{ code_dir }}
user=nobody
numprocs=1
stdout_logfile=/var/log/celery-worker.log
stderr_logfile=/var/log/celery-worker.log
autostart=true
autorestart=true

startsecs=10
stopwaitsecs = 600
killasgroup=true


[program:celerybeat]
command=/usr/local/bin/celery beat -A welgo:celery --schedule /var/lib/celery-beat.db --loglevel=INFO

directory={{ code_dir }}
user=root
numprocs=1
stdout_logfile=/var/log/celery-beat.log
stderr_logfile=/var/log/celery-beat.log
autostart=true
autorestart=true
startsecs=10

[group:celery]
programs=celery,celerybeat
